# conf/obmc-console-ssh@.service.in
[Unit]
Description=Console Client on %i
Wants=sshdgenkeys.service
After=sshdgenkeys.service

[Service]
# Unit defaults to a templated service
ExecStart=/usr/sbin/sshd -f /etc/ssh/sshd_config_port_%i -p %i
SyslogIdentifier=console-client-%i
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=default.target
